# supervisord.conf
[supervisord]
nodaemon=true
loglevel=DEBUG

[program:jupyter]
command=/start-jupyter.sh

#[program:tensorboard]
#command=tensorboard --logdir=/home/tensor/output
